# root/etc/crontab
*/2 * * * * /usr/sbin/run_sched.sh > /proc/1/fd/1 2>&1
